# ops/demo/supervisor.cfg
[supervisord]
nodaemon=true
pidfile=/tmp/supervisor.pid
logfile=/logs/supervisor.log
environment=LANG=C.UTF-8

[program:nginx]
command=nginx -c /ella/ops/demo/nginx.conf
redirect_stderr=true
stdout_logfile=/logs/nginx.log

[program:api]
command=gunicorn -c /ella/ops/demo/gunicorn.conf.py api.main:app
environment=DB_URL="postgresql:///postgres",PYTHONIOENCODING="utf-8",PYTHONUNBUFFERED="true",PYTHONPATH="/ella/src"
directory=/ella
redirect_stderr=true
stdout_logfile=/logs/api.log

[program:postgres]
command=/ella/ops/common/common_pg_startup
environment=PGDATA=/data
redirect_stderr=true
stdout_logfile=/logs/db.log

[program:gulp-production]
command=/ella/ops/demo/gulp
redirect_stderr=true
stdout_logfile=/logs/gulp.log
autorestart=false
