# ops/demo/supervisor.cfg
[supervisord]
nodaemon=true
pidfile=/socket/supervisor.pid
logfile=/logs/supervisor.log
environment=LANG=C.UTF-8
umask=002

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/socket/supervisor.sock
chmod=0770
username=ella_demo
password=ella_demo

[supervisorctl]
serverurl=unix:///socket/supervisor.sock
username=ella_demo
password=ella_demo

[program:nginx]
command=/ella/ops/common/common_nginx_startup
redirect_stderr=true
stdout_logfile=/logs/nginx.log

[program:api]
command=gunicorn -c /ella/ops/demo/gunicorn.conf.py api.main:app
environment=DB_URL="postgresql:///postgres",PYTHONIOENCODING="utf-8",PYTHONUNBUFFERED="true",PYTHONPATH="/ella/src"
directory=/ella
redirect_stderr=true
stdout_logfile=/logs/api.log

[program:postgres]
command=/ella/ops/common/common_pg_startup init
redirect_stderr=true
stdout_logfile=/logs/db.log