# ops/test/supervisor-export.cfg
[unix_http_server]
file = /tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[supervisord]
# in foreground to keep container running:
nodaemon=true
pidfile=/tmp/supervisor.pid
logfile=/var/log/supervisor.log
environment=LANG=C.UTF-8

[program:postgres]
command=/ella/ops/common/common_pg_startup
environment=PGDATA=/data
stopsignal=TERM
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

