# ops/test/supervisor.cfg
[supervisord]
nodaemon=false
pidfile=/tmp/supervisor.pid
logfile=/var/log/supervisor.log
environment=LANG=C.UTF-8

[program:postgres]
command=/ella/ops/pg_startup
environment=PGDATA=/data

[program:testannotationserver]
command=python /ella/ops/test/testannotationserver.py
stdout_logfile=foo.log
stdout_logfile_maxbytes=0
stdout_logfile=foo.log
stderr_logfile_maxbytes=0
redirect_stderr=true
