# ops/test/supervisor.cfg
[supervisord]
# in foreground to keep docker alive
nodaemon=true
pidfile=/tmp/supervisor.pid
logfile=/var/log/supervisor.log
environment=LANG=C.UTF-8

[program:postgres]
command=/ella/ops/common/common_pg_startup init
environment=PGDATA=/data

[program:testannotationserver]
command=python /ella/ops/test/testannotationserver.py
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
