# ops/test/supervisor.cfg
[supervisord]
nodaemon=false
pidfile=/tmp/supervisor.pid
logfile=/var/log/supervisor.log
environment=LANG=C.UTF-8

[program:postgres]
command=/ella/ops/pg_startup
environment=PGDATA=/data
