# ops/test/supervisor-e2e.cfg
[unix_http_server]
file = /tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[supervisord]
nodaemon=false
pidfile=/tmp/supervisor.pid
logfile=/dev/stdout
logfile_maxbytes=0
environment=LANG=C.UTF-8

# gulp started manually through supervisorctl after symlinking node_modules:
[program:gulp]
command=/dist/node_modules/gulp/bin/gulp.js
autostart=false
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:web]
command=python -u src/api/main.py
environment=SERVE_STATIC=true,USER_CONFIRMATION_ON_STATE_CHANGE="$(USER_CONFIRMATION_ON_STATE_CHANGE)s",USER_CONFIRMATION_TO_DISCARD_CHANGES="$(USER_CONFIRMATION_TO_DISCARD_CHANGES)s",DB_URL="postgresql:///postgres",PYTHONIOENCODING="utf-8",PYTHONUNBUFFERED="true",DEVELOP="true",ATTACHMENT_STORAGE="/ella/attachments"
#user=postgres
directory=/ella
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:postgres]
command=/ella/ops/common/common_pg_startup
environment=PGDATA=/data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
