# ops/prod/supervisor.cfg
[supervisord]
nodaemon=true
pidfile=/socket/supervisor.pid
logfile=/logs/supervisor.log
environment=LANG=C.UTF-8

[supervisorctl]
serverurl=unix:///socket/supervisor.sock

[program:nginx]
command=/ella/ops/common/common_nginx_startup
redirect_stderr=true
stdout_logfile=/logs/nginx.log
stdout_logfile_maxbytes=0

[program:api]
command=gunicorn -c /ella/ops/prod/gunicorn.conf.py api.main:app
environment=PYTHONIOENCODING="utf-8",PYTHONUNBUFFERED="true",PYTHONPATH="/ella/src"
directory=/ella
redirect_stderr=true
stdout_logfile=/logs/api.log
stdout_logfile_maxbytes=512MB
stdout_logfile_backups=99999999

[program:polling]
command=python src/api/polling.py
environment=PYTHONIOENCODING="utf-8",PYTHONUNBUFFERED="true",PYTHONPATH="/ella/src"
directory=/ella
redirect_stderr=true
stdout_logfile=/logs/polling.log
stdout_logfile_maxbytes=512MB
stdout_logfile_backups=99999999

[program:analysis-watcher]
environment=PYTHONIOENCODING="utf-8",PYTHONUNBUFFERED="true",PYTHONPATH="/ella/src"
command=python /ella/src/vardb/watcher/main.py --analyses $ANALYSIS_INCOMING --dest $ANALYSES_PATH --genepanels /ella/src/vardb/testdata/clinicalGenePanels/
redirect_stderr=true
stdout_logfile=/logs/analysis-watcher.log
stdout_logfile_maxbytes=0
